# images/fcos-base-image/root/usr/lib/systemd/system/cache-clean-ipi.timer
[Unit]
Description=Timer to clean cached iso unused and older than 15 days

[Timer]
OnCalendar=*-*-* *:30:00/3

[Install]
WantedBy=timers.target
