# images/fcos-bastion-image/root/usr/lib/systemd/system/nodes-pruner.timer
[Unit]
Description=Timer to prune cluster nodes running longer than 24 hour

[Timer]
OnCalendar=*-*-* *:02:00

[Install]
WantedBy=timers.target
